# contrib/tcp-intercept.conf
# Upstart script for tcp-intercept
description	"TCP intercept daemon"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

script
	# defaults
	LISTEN_PORT=5000
	LOGFILE="/var/log/tcp-intercept.log"

	if [ -r /etc/default/tcp-intercept ]; then
		. /etc/default/tcp-intercept
	fi

	/usr/local/bin/tcp-intercept >> "$LOGFILE"
end script
